# systemd/aivault.service
[Unit]
Description=AIVault - AI-powered knowledge vault
After=network.target

[Service]
Type=simple
User=aivault
WorkingDirectory=/opt/aivault
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=/opt/aivault/.env
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
